# conf/supervisor.conf
[supervisord]
nodaemon = true

[supervisorctl]
serverurl = unix:///run/supervisord.sock

[unix_http_server]
file = /run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command = nginx -g "daemon off;"

[program:php-server1]
command = php -S localhost:10001 demo/s1 

[program:php-server2]
command = php -S localhost:10002 demo/s2 

[program:php-server3]
command = php -S localhost:10003 demo/s3 

[program:sleepy-php]
command = php -S localhost:10004 demo/slow.php 

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
